# timetracker.service
[Unit]
Description=Timetracker Collector Daemon

[Service]
ExecStart=/TimeTracker/collect.py

[Install]
WantedBy=default.target
